# raskladka.service
[Unit]
Description=raskladka - keyboard layout switcher (double Shift)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/raskladka
Restart=on-failure
RestartSec=2
# For Wayland: add Environment=WAYLAND_DISPLAY, XDG_RUNTIME_DIR, YDOTOOL_SOCKET, DBUS_SESSION_BUS_ADDRESS
# Example: systemctl --user edit raskladka

[Install]
WantedBy=default.target
