# raskladka.service
[Unit]
Description=raskladka - keyboard layout switcher (double Shift)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/raskladka
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
