# raskladka.service
[Unit]
Description=raskladka - keyboard layout switcher (double Shift)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/Загрузки/raskladka/target/release/raskladka
Restart=on-failure
RestartSec=2
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
